# udev/96-battery.rules
KERNEL=="BAT0", SUBSYSTEM=="power_supply", RUN+="/usr/bin/pkill -SIGRTMIN+7 dwmblocks"
